# systemd/wide-dhcpv6-client.service.d/iserv-overrides.conf
[Service]
ExecReload=/bin/kill -SIGHUP $MAINPID
PIDFile=/run/dhcp6c.pid
